# systemd/user/stock-advisor-feishu-bot.service
[Unit]
Description=Stock Advisor Feishu Bot Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/projects/stock-advisor-bot
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 -m stock_advisor.cli serve-feishu-bot --config %h/projects/stock-advisor-bot/config.yaml
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
